# supervisor.docker.sft.conf
[program:supervisor.docker.sft]
command=make rebuilt
directory=/home/sft
autostart=true
autorestart=false